# streamer_proxy/systemd/streamer-proxy.service
[Unit]
Description=Dota Twitch Streamer Proxy
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/dota-twitch-lobby-bot/streamer_proxy
EnvironmentFile=/opt/dota-twitch-lobby-bot/streamer_proxy/.env
ExecStart=/opt/dota-twitch-lobby-bot/streamer_proxy/.venv/bin/uvicorn app.main:app --host ${APP_HOST} --port ${APP_PORT}
Restart=always
RestartSec=5
User=dota-streamer-proxy
Group=dota-streamer-proxy
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
